# chef/cookbooks/main/templates/default/init.d/unicorn.erb
#!/bin/bash
#
# unicorn        Unicorn app server
#
# chkconfig: - 85 15
# description: unicorn
# processname: unicorn
# pidfile: <%= @unicorn_pid %>

pidfile=<%= @unicorn_pid %>

start() {
  RACK_ENV=<%= @unicorn_environment %> unicorn -c <%= @unicorn_config %> -D
  echo 'done'
  RETVAL=0
  return $RETVAL
}

stop() {
  kill `cat <%= @unicorn_pid %>`
  RETVAL=0
  return $RETVAL
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 2
    start
    ;;
  *)
    echo $"Usage: $prog {start|stop}"
    exit 1
esac

exit $RETVAL

